# 99-qubes-tor-hook.rules
SUBSYSTEMS=="xen", KERNEL=="eth*", ACTION=="add", RUN+="/usr/lib/qubes-tor/start_tor_proxy.sh"
